# templates/systemd/docker-backup.timer.j2
#jinja2: trim_blocks: True, lstrip_blocks: True
{{ ansible_managed | comment }}
[Unit]
Description=Run Docker backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
